# qbittorrent-wireguard-sync.service
[Unit]
# Template schema version - increment when service structure changes
X-Version=3
Description=qBittorrent WireGuard Port Sync Service
Documentation=https://github.com/Simon-CR/qbittorrent-wireguard-pmp
After=network-online.target qbittorrent.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=%h/qbittorrent-wireguard-pmp
ExecStart=%h/qbittorrent-wireguard-pmp/port-sync.sh --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qbittorrent-wireguard-sync

# Security settings
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=%h/qbittorrent-wireguard-pmp

# Environment
Environment=WG_INTERFACE=wg0
Environment=QB_PORT=8080
Environment=NATPMP_GATEWAY=

# Notes:
# - The script runs in daemon mode, checking for port changes every 30 seconds.
# - All detection, connection, and update events are logged to port-sync.log and systemd journal.
# - For troubleshooting, check journalctl -u qbittorrent-wireguard-sync and port-sync.log.

[Install]
WantedBy=multi-user.target