# qbittorrent-wireguard-sync.service
[Unit]
Description=qBittorrent WireGuard Port Sync Service
Documentation=https://github.com/Simon-CR/qbittorrent-wireguard-pmp
After=network.target qbittorrent.service wireguard-quick@%i.service
Wants=network.target
BindsTo=wireguard-quick@%i.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=%h/qbittorrent-wireguard-pmp
ExecStart=%h/qbittorrent-wireguard-pmp/port-sync.sh --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qbittorrent-wireguard-sync

# Security settings
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=%h/qbittorrent-wireguard-pmp

# Environment
Environment=WG_INTERFACE=wg0
Environment=QB_PORT=8080

[Install]
WantedBy=multi-user.target